# datasciencetoolbox/sonartest/supervisord.conf
[supervisord]
nodaemon=true

[program:sonarqube]
command=./bin/run.sh #To adapt, I don't know how to launch your mongodb process

[program:sonarscanner]
command=sonar-scanner -DFOREGROUND